# backend/deploy/supervisor-biodata.conf
; Supervisor program: keeps the Gunicorn process alive and restarts it
; automatically on crash.  Drop this file at /etc/supervisor/conf.d/biodata.conf
; then run:
;
;     sudo supervisorctl reread && sudo supervisorctl update
;
[program:biodata]
command=/opt/biodata/venv/bin/gunicorn
        -c /opt/biodata/backend/deploy/gunicorn.conf.py
        biodata_system.wsgi:application
directory=/opt/biodata/backend
user=biodata
group=biodata
umask=007

; Load Django settings from .env (python-decouple reads from process env or
; the .env file in the project dir; we let python-decouple find it).
environment=PYTHONUNBUFFERED="1",DJANGO_SETTINGS_MODULE="biodata_system.settings"

; --- Lifecycle --------------------------------------------------------------
autostart=true
autorestart=true
startretries=5
startsecs=5
stopsignal=TERM
stopwaitsecs=30

; --- Logging ----------------------------------------------------------------
stdout_logfile=/var/log/biodata/gunicorn.out.log
stderr_logfile=/var/log/biodata/gunicorn.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_backups=5

; --- Sanity -----------------------------------------------------------------
; If you need to wipe socket dir before start (e.g. /run is tmpfs), prefer a
; systemd tmpfiles.d rule over a shell wrapper here.
